# file_mover_portal/examples/linux/file-mover.service
[Unit]
Description=Secure File Mover Portal
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=filemover
Group=filemover
WorkingDirectory=/opt/file_mover_portal
EnvironmentFile=/opt/file_mover_portal/.env
UMask=0077
ExecStart=/opt/file_mover_portal/.venv/bin/gunicorn --config /opt/file_mover_portal/examples/linux/gunicorn.conf.py app:app
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
ReadWritePaths=/opt/file_mover_portal/runtime /opt/file_mover_portal/logs /data/incoming /data/processed

[Install]
WantedBy=multi-user.target
